# systemd/netspecter-speedtest.service
[Unit]
Description=NetSpecter Scheduled Internet Speed Test
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=netspecter
Group=netspecter
ExecStart=/opt/netspecter/venv/bin/python /opt/netspecter/scheduled_speedtest.py
